# deploy/systemd/investment-agent-migrate.service
[Unit]
Description=Investment Agent Database Migration
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=deploy_invest_ai
Group=deploy_invest_ai
WorkingDirectory=/home/deploy_invest_ai/investment-agent
Environment=PATH=/home/deploy_invest_ai/.local/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/home/deploy_invest_ai/investment-agent/.env
ExecStart=/home/deploy_invest_ai/investment-agent/scripts/run_migrations.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=investment-agent-migrate

[Install]
WantedBy=multi-user.target
